# templates/upstart.conf.j2
start on login-session-start

script

/usr/bin/x11vnc {{vnc_params}}

end script
